# run/service/s6-svscan.log/run
#!/bin/rc -e

for (logd in /run/uncaught-logs /var/log/init) {
	[ -d $logd ] || install -d -o catchlog -g adm -m 0750 $logd
}

exec \
redirfd -w 1 /dev/console \
fdmove -c 2 1 \
redirfd -rnb 0 fifo \
s6-setuidgid catchlog \
emptyenv -p \
s6-log -bp t T \
	!'/etc/slew/init/save_log.rc /var/log/init 6 5000000' \
	/run/uncaught-logs

