# run/service/s6-svscan.log/run
#!/bin/rc -e

exec > /dev/console >[2=1]
[ -d /run/uncaught-logs ] ||
	install -d -o catchlog -g adm -m 0750 /run/uncaught-logs

exec \
redirfd -rnb 0 fifo \
s6-setuidgid catchlog \
emptyenv -p \
s6-log -bp t T \
	!'/etc/slew/init/save_log.rc /var/log/init 6 5000000' \
	/run/uncaught-logs

